# packages/daemon/systemd/comis.service
# Comis Daemon - systemd unit file
#
# Security-hardened service definition for production deployment.
# Requires Node.js 24+ for --permission flag support.
#
# Install:
#   cp comis.service /etc/systemd/system/
#   systemctl daemon-reload
#   systemctl enable --now comis
#
# SVC-01: systemd service management
# SVC-02: Security hardening (ProtectSystem, MemoryDenyWriteExecute, SystemCallFilter)
# SVC-03: Watchdog health monitoring (WatchdogSec + sd-notify)

[Unit]
Description=Comis AI Agent Daemon
Documentation=https://github.com/comis/comis
After=network-online.target
Wants=network-online.target

[Service]
Type=notify
User=comis
Group=comis
WorkingDirectory=/opt/comis

# --jitless: Disable JIT to reduce attack surface (no W^X pages)
# --permission: Enable Node.js permission model for filesystem/network restrictions
ExecStart=/usr/bin/node --jitless --permission --allow-fs-read=/opt/comis --allow-fs-write=/var/lib/comis --allow-child-process dist/daemon.js

# Restart policy
Restart=on-failure
RestartSec=5s

# Watchdog: systemd kills process if no ping within WatchdogSec
# Daemon pings at half this interval (15s) with health gating
WatchdogSec=30s

# Shutdown: allow 45s for graceful shutdown (app timeout is 30s)
TimeoutStopSec=45

# Resource limits
MemoryMax=2G
TasksMax=100

# Standard output/error to journald for structured log collection
StandardOutput=journal
StandardError=journal
SyslogIdentifier=comis

# Environment
Environment=NODE_ENV=production
EnvironmentFile=-/etc/comis/env

# --- Security Hardening ---

# Filesystem protection
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
ReadWritePaths=/var/lib/comis /var/log/comis

# Memory protection
MemoryDenyWriteExecute=yes

# Privilege escalation prevention
NoNewPrivileges=yes

# Drop all capabilities
CapabilityBoundingSet=

# System call filtering: allow only standard service calls
SystemCallFilter=@system-service
SystemCallArchitectures=native

# Network namespace (allow outbound for LLM API calls)
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
ProtectClock=yes
ProtectHostname=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
RestrictNamespaces=yes
LockPersonality=yes

[Install]
WantedBy=multi-user.target
